# data/greeter-autostart/gok.desktop.in
[Desktop Entry]
_Name=GNOME Onscreen Keyboard
_Comment=Use an onscreen keyboard
TryExec=gok
Exec=gok --login --access-method=directselection
Terminal=false
Type=Application
StartupNotify=true
Categories=GNOME;GTK;Accessibility;
AutostartCondition=GNOME /desktop/gnome/applications/at/screen_keyboard_enabled
X-GNOME-Provides=screen-keyboard
X-GNOME-Bugzilla-Bugzilla=GNOME
X-GNOME-Bugzilla-Product=gok
X-GNOME-Bugzilla-Component=general

