# data/autostart/caribou-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=Caribou
Exec=@LIBEXECDIR@/caribou
AutostartCondition=GSettings org.gnome.desktop.a11y.applications screen-keyboard-enabled
X-GNOME-AutoRestart=true
