# data/61-gdm.rules.in
# disable Wayland on Hi1710 chipsets
ATTR{vendor}=="0x19e5", ATTR{device}=="0x1711", RUN+="@libexecdir@/gdm-runtime-config set daemon WaylandEnable false"
# default to Xorg on single gpu vendor nvidia systems
DRIVER=="nvidia", RUN+="@libexecdir@/gdm-runtime-config set daemon PreferredDisplayServer xorg"
# disable Wayland on hybrid systems with vendor nvidia driver
SUBSYSTEM=="drm", KERNEL=="card[1-9]*", TEST=="/sys/module/nvidia_drm", RUN+="@libexecdir@/gdm-runtime-config set daemon WaylandEnable false"
# disable Wayland on systems with vendor nvidia driver but modesetting disabled
ACTION=="add", SUBSYSTEM=="module", KERNEL=="nvidia_drm", ATTR{parameters/modeset}=="N", RUN+="@libexecdir@/gdm-runtime-config set daemon WaylandEnable false"
# disable Wayland if modesetting is disabled
IMPORT{cmdline}="nomodeset", RUN+="@libexecdir@/gdm-runtime-config set daemon WaylandEnable false"
