# data/gdm.service.in
[Unit]
Description=GNOME Display Manager

# replaces the getty
Conflicts=getty@tty@GDM_INITIAL_VT@.service
After=getty@tty@GDM_INITIAL_VT@.service

# replaces plymouth-quit since it quits plymouth on its own
Conflicts=plymouth-quit.service
After=plymouth-quit.service

# Needs all the dependencies of the services it's replacing
# pulled from getty@.service and plymouth-quit.service
# (except for plymouth-quit-wait.service since it waits until
# plymouth is quit, which we do)
After=rc-local.service plymouth-start.service systemd-user-sessions.service

# GDM takes responsibility for stopping plymouth, so if it fails
# for any reason, make sure plymouth still stops
OnFailure=plymouth-quit.service

[Service]
ExecStart=@sbindir@/gdm
KillMode=process
Restart=always
IgnoreSIGPIPE=no
BusName=org.gnome.DisplayManager
StandardOutput=syslog
StandardError=inherit
EnvironmentFile=-@LANG_CONFIG_FILE@

[Install]
Alias=display-manager.service
