# trunk/debian/postrm
#!/bin/sh
#
# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
	update-rc.d gdm remove >/dev/null
fi
# End automatically added section
#DEBHELPER#
